# roles/hue/templates/hue.service.j2
[Unit]
Description=Hue

[Service]
User={{ hue_user }}
Type=simple

ExecStart={{ hue_service_start_script }}
ExecStop={{ hue_service_stop_script }}

[Install]
WantedBy=multi-user.target
